# data/rules/hplip-printer@.service
[Unit]
Description=Configure HP Plugged-In Printer

[Service]
ExecStart=/usr/bin/python /usr/bin/hp-config_usb_printer "%i"
